# build/copyfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=service squid start
priority=10

[program:ipsec]
command=ipsec start
priority=20

[program:xl2tpd]
command=service xl2tpd start
priority=30

[program:cron]
command=service cron start
priority=40

[program:vpn-routes-add]
#command=/root/.config_files/vpn-net-route-add.sh
command=sh /root/.config_files/vpn-net-route-add.sh
priority=50
startsecs = 0
autorestart=false
startretries = 1
